# config/imperium.service
[Unit]
Description=Imperium Intent-Based Networking Controller
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Imperium
Environment="PATH=/home/pi/Imperium/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_ENV=production"
Environment="JWT_SECRET_KEY=CHANGE_THIS_IN_PRODUCTION"
Environment="NETWORK_INTERFACE=eth0"

# Ensure database directory exists
ExecStartPre=/bin/mkdir -p /home/pi/Imperium/data
ExecStartPre=/bin/chown pi:pi /home/pi/Imperium/data

# Start the Imperium controller
ExecStart=/home/pi/Imperium/venv/bin/python3 /home/pi/Imperium/src/main.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=false
PrivateTmp=true

# Logging
StandardOutput=append:/var/log/imperium/controller.log
StandardError=append:/var/log/imperium/controller-error.log

# Resource limits
MemoryLimit=512M
CPUQuota=60%

[Install]
WantedBy=multi-user.target
